# systemd/telegram-bot.service
[Unit]
Description=Telegram Bot with Ollama
After=network.target ollama.service
Requires=ollama.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/TakeYourPillsBot
Environment="PATH=/root/TakeYourPillsBot/venv/bin"
EnvironmentFile=/root/TakeYourPillsBot/.env
ExecStart=/root/TakeYourPillsBot/venv/bin/python /root/TakeYourPillsBot/bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target